# installer_files/conda/pkgs/conda-25.3.0-py312h2e8e312_0/etc/profile.d/conda.csh
setenv CONDA_EXE "`cygpath 'D:\bld\conda_1743046486868\_h_env\Scripts\conda.exe'`";
setenv _CONDA_ROOT "`cygpath 'D:\bld\conda_1743046486868\_h_env'`";
setenv _CONDA_EXE "`cygpath 'D:\bld\conda_1743046486868\_h_env\Scripts\conda.exe'`";
setenv CONDA_PYTHON_EXE "`cygpath 'D:\bld\conda_1743046486868\_h_env\python.exe'`";
source "`cygpath 'D:\bld\conda_1743046486868\_h_env\Lib\site-packages\conda\shell\etc\profile.d\conda.csh'`";
