# installer_files/conda/etc/profile.d/conda.csh
setenv CONDA_EXE "`cygpath 'C:/Orion/text-generation-webui/installer_files/conda\Scripts\conda.exe'`";
setenv _CONDA_ROOT "`cygpath 'C:/Orion/text-generation-webui/installer_files/conda'`";
setenv _CONDA_EXE "`cygpath 'C:/Orion/text-generation-webui/installer_files/conda\Scripts\conda.exe'`";
setenv CONDA_PYTHON_EXE "`cygpath 'C:/Orion/text-generation-webui/installer_files/conda\python.exe'`";
source "`cygpath 'C:/Orion/text-generation-webui/installer_files/conda\Lib\site-packages\conda\shell\etc\profile.d\conda.csh'`";
